# home-cluster/00-init.files/schedule-reboot.timer
[Unit]
Description=Reboot Scheduling.

[Timer]
OnCalendar=Sat *-*-* 04:00:00
Unit=reboot.target

[Install]
WantedBy=multi-user.target
